# systemd/alpr-rsync.timer
[Unit]
Description=Sync feeds every 15 minutes

[Timer]
OnBootSec=3min
OnUnitInactiveSec=900

[Install]
WantedBy=timers.target
